# stage7/02-setup-wlan/files/create_ap_interface.service
[Unit]
Description=Create Virtual WLAN Interface
Wants=network-pre.target
Before=network-pre.target
BindsTo=sys-subsystem-net-devices-wlan0.device
After=sys-subsystem-net-devices-wlan0.device

[Service]
Type=oneshot
ExecStart=/usr/local/bin/create_ap_interface.sh

[Install]
WantedBy=multi-user.target
